# recipes-core/util-linux/util-linux/fstrim.service
[Unit]
Description=Discard unused blocks
ConditionVirtualization=!container

[Service]
Type=oneshot
ExecStart=@SBINDIR@/fstrim -av
